# systemd/fyzzy-bridge.service
[Unit]
Description=Fyzzy Bridge (Keiser Hub collector)
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=fyzzy
WorkingDirectory=/opt/fyzzy-bridge
EnvironmentFile=-/etc/fyzzy-bridge.env
# Runs the single bundled file behind the `current` symlink, which OTA swaps.
ExecStart=/usr/bin/node /opt/fyzzy-bridge/current/bundle.cjs
Restart=always
RestartSec=5
# Provisioning needs to run nmcli; keep it simple for the appliance.
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
